# rpi/wyoming-satellite/wyoming-satellite.service
[Unit]
Description=Grumpy Wyoming Satellite
# Only meaningful once networking is up — we forward mic audio to the remote
# wake word service and HA reaches us over the network.
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
# User/WorkingDirectory/EnvironmentFile are filled in by setup.sh when it
# copies this template into /etc/systemd/system/.
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/config.env
# PipeWire's socket lives under the user's runtime dir. Needed so pw-record /
# pw-play (used below) can find the PipeWire daemon from this system service.
Environment=XDG_RUNTIME_DIR=/run/user/__USER_UID__

# Satellite advertises itself on ${SATELLITE_PORT} for HA to connect to.
# Mic and speaker go through PipeWire (pw-record / pw-play) rather than direct
# ALSA (arecord / aplay) so the Chromium kiosk and the satellite can share the
# same USB mic/speaker — direct ALSA locks the hardware exclusively.
# --vad gates audio streaming so we only send speech to STT, not silence.
ExecStart=__INSTALL_DIR__/.venv/bin/python -m wyoming_satellite \
    --name ${DEVICE_NAME} \
    --uri tcp://0.0.0.0:${SATELLITE_PORT} \
    --mic-command "pw-record --rate=${MIC_RATE} --channels=1 --format=s16 --raw -" \
    --snd-command "pw-play --rate=${SPEAKER_RATE} --channels=1 --format=s16 --raw -" \
    --wake-uri tcp://${WAKE_HOST}:${WAKE_PORT} \
    --wake-word-name ${WAKE_WORD_NAME} \
    --mic-auto-gain 5 \
    --mic-noise-suppression 1 \
    --vad

Restart=on-failure
RestartSec=3s

# Low-level audio needs these; harmless on systems that don't enforce them.
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
