# rpi/dashboard/grumpy-dashboard.service
[Unit]
Description=Grumpy Dashboard Dev Server
After=network.target

[Service]
Type=simple
User=jason
WorkingDirectory=/home/jason/grumpy/rpi/dashboard
ExecStart=/usr/bin/npx vite --host
Restart=always
RestartSec=5
Environment=NODE_ENV=development
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
